# net-misc/ibrdtnd/files/dtnd.service
[Unit]
Description=dtnd

[Service]
ExecStart=/usr/sbin/dtnd -c /etc/ibrdtnd.conf

[Install]
WantedBy=multi-user.target
